# docker/supervisor/elasticsearch.conf
[program:elasticsearch]
command=/usr/share/elasticsearch/bin/elasticsearch
autostart=true
autorestart=true
stdout_logfile=/var/log/elasticsearch.log
redirect_stderr=true
user=root
