# docker/supervisor/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -c /etc/kibana.conf
autostart=true
autorestart=true
user=root
